# scheduler/memory-gardener.service
[Unit]
Description=Claude Code Memory Gardener

[Service]
Type=oneshot
WorkingDirectory=__REPO__
ExecStart=__PYTHON__ -m gardener run
